# plexmediaserver.conf
[program:plexmediaserver]
user = nobody
command = /scripts/start.sh
